# other-config/systemd/user/nobeep.service
[Unit]
Description=stop pcspkr module

[Service]
ExecStart=bash -c "echo 'peach' | sudo -S modprobe -r pcspkr"

[Install]
WantedBy=default.target
